# orvix-main/release/systemd/orvix.service
[Unit]
Description=Orvix Email Server Platform
Documentation=https://orvix.email/docs
After=network.target
Wants=network.target

[Service]
Type=simple
User=orvix
Group=orvix
ExecStart=/usr/local/bin/orvix serve
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10

# Security hardening
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true

# File limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=ORVIX_CONFIG=/etc/orvix/orvix.yaml
Environment=ORVIX_LOG_DIR=/var/log/orvix

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orvix

[Install]
WantedBy=multi-user.target

# RC2 FIX: Removed StartLimitBurst and StartLimitIntervalSec
# These directives are invalid on Ubuntu 22.04+ systemd
# The Restart=always and RestartSec=10 provide sufficient restart behavior